# packages/worker/systemd/autoapply-worker@.service
[Unit]
Description=AutoApply Worker %i
After=network.target

[Service]
Type=simple
User=autoapply
WorkingDirectory=/opt/autoapply/worker
Environment=WORKER_ID=worker-%i
EnvironmentFile=/opt/autoapply/.env
ExecStart=/opt/autoapply/venv/bin/python worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
